# etc/profile.d/go_env
GOROOT=$HOME/opt/gobin

GOPATH=$HOME/opt/goapp

PATH=$PATH:$GOROOT/bin

export GOROOT PATH
