# etc/profile.d/go_env
GOROOT=$HOME/opt/go

PATH=$PATH:$GOROOT/bin

export GOROOT PATH
