# etc/profile.d/gradle_env
GRADLE_HOME=$HOME/opt/gradle

PATH=$PATH:$GRADLE_HOME/bin

export PATH
