# etc/profile.d/kafka_env
KAFKA_HOME=$HOME/opt/kafka

PATH=$PATH:$KAFKA_HOME/bin

export PATH KAFKA_HOME
