# etc/profile.d/sbt_env
SBT_HOME=$HOME/opt/sbt

PATH=$SBT_HOME/bin:$PATH

export PATH SBT_HOME
