# etc/profile.d/scala_env
SCALA_HOME=$HOME/opt/scala

PATH=$SCALA_HOME/bin:$PATH

export PATH
