# etc/profile.d/spark_env
SPARK_HOME=$HOME/opt/spark

PATH=$PATH:$SPARK_HOME/bin

export PATH SPARK_HOME
